# targets/hiker-x9/package/network/services/hiker-x9-standard-defaults/files/hiker-boot-metrics.init
#!/bin/sh /etc/rc.common

START=02
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /bin/sh -c 'logger -t hiker-boot "rc start +$(cut -d. -f1 /proc/uptime)s"'
	procd_set_param respawn 0 0 0
	procd_close_instance
}

